# cosmos-node/templates/cosmos.service.j2
# Managed by Ansible - DO NOT EDIT MANUALLY
# {{ daemon_name }} systemd service

[Unit]
Description={{ chain | upper }} Blockchain Node ({{ daemon_name }})
After=network-online.target
Wants=network-online.target

[Service]
User={{ daemon_user }}
Group={{ daemon_user }}
ExecStart=/usr/local/bin/cosmovisor run start --home {{ daemon_home }}
Restart=always
RestartSec=3
LimitNOFILE=1048576
LimitNPROC=65535

# Cosmovisor environment
Environment="DAEMON_HOME={{ daemon_home }}"
Environment="DAEMON_NAME={{ daemon_name }}"
Environment="DAEMON_ALLOW_DOWNLOAD_BINARIES={{ cosmovisor_allow_download | lower }}"
Environment="DAEMON_RESTART_AFTER_UPGRADE={{ cosmovisor_restart_after_upgrade | lower }}"
Environment="DAEMON_LOG_BUFFER_SIZE=512"
Environment="UNSAFE_SKIP_BACKUP=true"

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier={{ daemon_name }}

[Install]
WantedBy=multi-user.target
